# deploy/overlay-api.service
# Install:
#   sudo install -m 644 deploy/overlay-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now overlay-api.service
#
# Adjust User, Group, WorkingDirectory, and ExecStart to match your install (e.g. /opt/fastapi-app).

[Unit]
Description=Logo overlay FastAPI (uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=overlay
Group=overlay
WorkingDirectory=/opt/fastapi-app
Environment=PATH=/opt/fastapi-app/.venv/bin:/usr/local/bin:/usr/bin
# Use python -m uvicorn (not the uvicorn console script) to avoid systemd 203/EXEC
# when the wrapper shebang or SELinux blocks the script.
ExecStart=/opt/fastapi-app/.venv/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=overlay-api

[Install]
WantedBy=multi-user.target
